# etc/supervisord.d/hdfs.ini
[program:hdfs]
command=/opt/qnib/hadoop/bin/hdfs.sh
redirect_stderr=true
startsecs=0
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
