# etc/supervisord.d/kafka_update.ini
[program:kafka-update]
command=consul-template -consul localhost:8500 -template "/etc/consul-templates/kafka.server.properties.ctmpl:/opt/kafka/config/server.properties:supervisorctl restart kafka"
autostart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
