# etc/supervisord.d/kafka.ini
[program:kafka]
command=/opt/qnib/kafka/bin/start_kafka.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
