# etc/supervisord.d/kafka_topics.ini
[program:kafka_topics]
command=/opt/qnib/kafka/bin/create_topics.sh
redirect_stderr=true
startsecs=0
exitcodes=9
stdout_logfile=/var/log/supervisor/%(program_name)s.log
