# etc/supervisord.d/sshd.ini
[program:sshd]
command=/opt/qnib/bin/startup_sshd.sh
redirect_stderr=true
stdout_syslog=true
